# luci-app-oc-editor/root/etc/uci-defaults/oc-editor.sh
#!/bin/sh

cd /
[ ! -d /www/oc-editor/rootfs ] && ln -s / /www/oc-editor/rootfs
exit 0
